# systemd/abraxas-oracle-daily.service
[Unit]
Description=Abraxas Oracle v2 Daily Run
Documentation=https://github.com/scrimshawlife-ctrl/Abraxas
After=network.target

[Service]
Type=oneshot
User=abraxas
Group=abraxas
WorkingDirectory=/opt/abraxas
Environment="PYTHONPATH=/opt/abraxas"

# Execute daily Oracle run
ExecStart=/usr/bin/python3 -m abraxas.cli.oracle_daily_run \
    --config /opt/abraxas/config/oracle_daily.json \
    --output-dir /var/lib/abraxas/oracle_history \
    --intel-dir /var/lib/abraxas/intel \
    --report-dir /var/lib/abraxas/reports \
    --log-file /var/log/abraxas/oracle_daily.log \
    --verbose

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=abraxas-oracle-daily

# Resource limits
MemoryLimit=1G
CPUQuota=50%

# Restart policy
Restart=on-failure
RestartSec=5min

[Install]
WantedBy=multi-user.target
